# 99-net-audit.rules
# /etc/udev/rules.d/99-net-audit.rules
# Start QuickPeek when an Ethernet link appears or carrier becomes active.
#
# Works on physical Ethernet-style names such as eth*, enp*, ens*.
# Does not trigger on loopback, bridge, VLAN, or WLAN interfaces.

ACTION=="add", SUBSYSTEM=="net", \
  KERNEL!="lo", KERNEL!="br*", KERNEL!="vlan*", KERNEL!="wlan*", \
  RUN+="/bin/systemctl start net_audit.service"

ACTION=="change", SUBSYSTEM=="net", \
  KERNEL!="lo", KERNEL!="br*", KERNEL!="vlan*", KERNEL!="wlan*", \
  ATTR{carrier}=="1", \
  RUN+="/bin/systemctl start net_audit.service"
