# net_audit.service
[Unit]
Description=QuickPeek network audit
After=network-online.target
Wants=network-online.target
BindsTo=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/python3 /opt/quickpeek/net_audit.py --config /opt/quickpeek/config.ini
Restart=no
TimeoutStartSec=900
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quickpeek
NoNewPrivileges=no
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
